# debian/manuka-api.init.in
#!/bin/sh
### BEGIN INIT INFO
# Provides:          manuka-api
# Required-Start:    $network $local_fs $remote_fs $syslog
# Required-Stop:     $remote_fs
# Should-Start:      postgresql mysql keystone rabbitmq-server ntp
# Should-Stop:       postgresql mysql keystone rabbitmq-server ntp
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Manuka API
# Description:       Manuka API daemon
### END INIT INFO

DESC="OpenStack Manuka API (manuka-api)"
PROJECT_NAME=manuka
NAME=${PROJECT_NAME}-api
